# sage-daemon/launchd/com.sage.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sage.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/lyf/.sage/bin/sage-daemon</string>
        <string>--config</string>
        <string>/Users/lyf/.sage/config.toml</string>
        <string>--foreground</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/lyf/.sage/logs/sage.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/lyf/.sage/logs/sage.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/lyf</string>
        <key>USER</key>
        <string>lyf</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/lyf/.cargo/bin:/opt/homebrew/bin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>http_proxy</key>
        <string>http://127.0.0.1:7890</string>
        <key>https_proxy</key>
        <string>http://127.0.0.1:7890</string>
        <key>all_proxy</key>
        <string>socks5://127.0.0.1:7890</string>
    </dict>
</dict>
</plist>
